# crontab/crontab-php-nginx
5 0 */10 * *  bash /opt/archive.sh  >> /var/log/crontab/archive_php_nginx.log
